# supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/sbin/mysqld
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/mysql.err.log
stdout_logfile=/var/log/supervisor/mysql.out.log

[program:apache2]
command=/usr/bin/pidproxy /var/run/apache2/apache2.pid /bin/bash -c "source /etc/apache2/envvars && /usr/sbin/apache2ctl -DFOREGROUND"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/apache2.err.log
stdout_logfile=/var/log/supervisor/apache2.out.log

[program:bootstrap]
command=/scripts/bootstrap.sh
autostart=true
autorestart=false
startretries=0
#redirect_stderr=false
#redirect_stdout=false
stderr_logfile=/var/log/supervisor/bootstrap.err.log
stdout_logfile=/var/log/supervisor/bootstrap.out.log
